# contrib/sshfail2kml-1.3.3/debian/sshfail2kml.cron.d
#
# Regular cron jobs for the sshfail2kml package
#
0 4 * * *	root	[ -x /usr/bin/sshfail2kml ] && /usr/bin/sshfail2kml -q
